# setup/sousvide-setup.conf
start on local-filesystems

task

script
    modprobe w1-gpio
    modprobe w1-therm
    modprobe bcm2708_wdog
    bash /root/setup.sh
end script
